# fedora/41/files/etc/systemd/system/k3s.service.d/override.conf
[Service]
ExecStartPost=/bin/chown root:k3s /etc/rancher/k3s/k3s.yaml
ExecStartPost=/bin/chmod 640 /etc/rancher/k3s/k3s.yaml
